# scripts/da3.service.in
[Unit]
Description=Depth Anything 3 Depth Estimation Service
After=network.target

[Service]
Type=simple
User=@USER@
WorkingDirectory=@REL_ROOT@
# Foreground mode: start.sh exec's the binary, so systemd tracks the real
# process. Do NOT add --daemon here - Type=simple would see the fork parent
# exit and restart-loop the service.
ExecStart=@REL_ROOT@/scripts/start.sh --port 8090
Restart=on-failure
RestartSec=5

# Environment
Environment=DA3_MODEL_DIR=@REL_ROOT@/resources/nnmodels
Environment=DA3_WWW_DIR=@REL_ROOT@/resources

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=da3-service

[Install]
WantedBy=multi-user.target
